# src/atlas2/atlas2.timer.template
[Unit]
Description=Run ATLAS2 enchanced asset discovery service on a Systemd timer

[Timer]
OnCalendar=[OnCalendar TIME UNIT (e.g. *-*-* 2:00:00)]
Persistent=true

[Install]
WantedBy=timers.target

